# tools/beaglebone/pin-config.service
[Unit]
Description=Setup for BBB pins

[Service]
Type=simple
ExecStart=/bin/bash /home/debian/pin-config.sh

[Install]
WantedBy=multi-user.target